# backend/supervisor/celery.conf
[program:celery_worker]
command=/root/.venv/bin/celery -A celery_app worker --loglevel=info --concurrency=2
directory=/app/backend
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/var/log/supervisor/celery_worker.out.log
stderr_logfile=/var/log/supervisor/celery_worker.err.log
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:celery_beat]
command=/root/.venv/bin/celery -A celery_app beat --loglevel=info
directory=/app/backend
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/var/log/supervisor/celery_beat.out.log
stderr_logfile=/var/log/supervisor/celery_beat.err.log
stopwaitsecs=10
