# backend/supervisor/celery.conf
; Set MINT_STREET_ROOT in the supervisor environment to the repo root,
; and CELERY_BIN to the celery executable (e.g. <venv>/bin/celery).
[program:celery_worker]
command=%(ENV_CELERY_BIN)s -A celery_app worker --loglevel=info --concurrency=2
directory=%(ENV_MINT_STREET_ROOT)s/backend
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/var/log/supervisor/celery_worker.out.log
stderr_logfile=/var/log/supervisor/celery_worker.err.log
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:celery_beat]
command=%(ENV_CELERY_BIN)s -A celery_app beat --loglevel=info
directory=%(ENV_MINT_STREET_ROOT)s/backend
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/var/log/supervisor/celery_beat.out.log
stderr_logfile=/var/log/supervisor/celery_beat.err.log
stopwaitsecs=10
