# backend/supervisor/telegram_scheduler.conf
[program:telegram_scheduler]
command=/root/.venv/bin/python3 /app/backend/telegram_scheduler.py
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/telegram.out.log
stderr_logfile=/var/log/supervisor/telegram.err.log
