# backend/supervisor/telegram_scheduler.conf
; Set MINT_STREET_ROOT in the supervisor environment to the repo root,
; and PYTHON_BIN to the Python interpreter (e.g. the venv python).
[program:telegram_scheduler]
command=%(ENV_PYTHON_BIN)s %(ENV_MINT_STREET_ROOT)s/backend/telegram_scheduler.py
directory=%(ENV_MINT_STREET_ROOT)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/telegram.out.log
stderr_logfile=/var/log/supervisor/telegram.err.log
